# deploy/personblog-sync.service
[Unit]
Description=Sync personblog from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/personblog
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash /var/www/personblog/deploy/sync-from-github.sh
